# systemd/backup-daily.service
[Unit]
Description=Backup quotidien automatique

[Service]
ExecStart=/opt/backup_daily.sh
